# .config/systemd/user/mbsync@.service
[Unit]

[Service]
ExecStart=/usr/bin/mbsync %i
