# .config/systemd/user/notmuch.timer
[Unit]
Description=Notmuch indexer timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m

[Install]
WantedBy=timers.target
